# roles/dn42_autopeer/templates/dn42-registry-pull.service.j2
[Unit]
Description=Update Git Repository

[Service]
Type=oneshot
ExecStart=/usr/bin/git -C {{ dn42_autopeer_repo_directory }} pull
User=root
